# node-agent/templates/linx-node-agent.service.example
[Unit]
Description=LinX Node Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=LINX_EXTERNAL_AGENT_COMMAND=external-agent --prompt-file "$LINX_AGENT_PROMPT_FILE"
ExecStart=%h/.local/bin/linx-node-agent \
  --base-url http://YOUR_SERVER:8000/api/v1 \
  --token YOUR_BEARER_TOKEN \
  --project-id YOUR_PROJECT_ID \
  --name "%H" \
  --capability host_execution \
  --capability shell \
  --capability ops
Restart=always
RestartSec=5
WorkingDirectory=%h

[Install]
WantedBy=default.target
